# vc_dataroom/com.amphoraxe.vc.webapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amphoraxe.vc.webapp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/amphoraxe/Developer/vc_dataroom/venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8010</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/amphoraxe/Developer/vc_dataroom</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/amphoraxe/Developer/vc_dataroom/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>VC_SECRET_KEY</key>
        <string>CHANGE_THIS_TO_A_SECURE_RANDOM_KEY</string>
        <key>VC_ADMIN_PASSWORD</key>
        <string>CHANGE_THIS_ADMIN_PASSWORD</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/amphoraxe/Developer/vc_dataroom/logs/webapp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amphoraxe/Developer/vc_dataroom/logs/webapp.error.log</string>
</dict>
</plist>
